# packaging/vconf-setup.service
[Unit]
Description=Initialize the VConf storage
DefaultDependencies=no
<<<<<<< HEAD
=======
RequiresMountsFor=/var/run
>>>>>>> 57306cd... Improve path dependency
Before=basic.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/etc/rc.d/init.d/vconf-init

[Install]
WantedBy=basic.target

